# config/supervisor_ffwd.ini
[supervisord]
nodaemon=true

[group:djangoes]
programs=django-8001,django-8002

[program:django-8001]
command=python /home/django/ffwdsite/tornado_run.py 8001
directory=/home/django/ffwdsite
autorestart=true
redirect_stderr=true
loglevel=info
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:django-8002]
command=python /home/django/ffwdsite/tornado_run.py 8002
directory=/home/django/ffwdsite
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tornado.log
loglevel=info

[program:nginx]
command=nginx
directory=/usr/sbin
autorestart=true
redirect_stderr=true
stdout_logfile=/tmp/nginx.error
loglevel=information