# config/supervior_celery.ini
[program:celery]
command=celery -A tasks.celery_task.app worker -B -l warn -c 1
directory=/home/django/
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
